# installers/rhel8/dnsblock.service
[Unit]
Description=update the DNS blocklist

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/dnsblock
